# crontab
05	*	*	*	*	root	/usr/bin/LinkedinToTelegram/LinkedinToTelegram.sh
